# .crontab
# Define o PATH
PATH=/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin

# Cronjobs
# Brew
10 3 * * * brew list > ~/.brew/applications.txt
0 3 * * 1,6 brew update && brew upgrade

# crontab
15 3 * * * crontab -l > ~/.crontab

# SDKMAN
0 5 * * 6 ~/.scripts/update_java.sh

# YADM
30 3 * * * cd ~ && (yadm status --porcelain | grep '^[ M| D]' && yadm add -u && yadm commit -m "dotfiles automatic update") && yadm push -u origin main
