# .crontab
# Define o PATH
PATH=/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin

# Cronjobs
# Brew
10 3 * * * brew list > ~/.brew/applications.txt
0 3 * * 1,6 brew update && brew upgrade

# crontab
15 3 * * * crontab -l > ~/.crontab

# YADM
30 3 * * * cd ~ && [ -n "$(yadm status --porcelain)" ] && yadm add -u && yadm commit -am "dotfiles automatic update" && yadm push -u origin main
